# root/etc/services.d/ssh/run
#!/usr/bin/with-contenv sh

echo "Starting ssh server..."
exec /usr/sbin/sshd -D -e
